# customization/includes.chroot/etc/systemd/system/ssh.service.d/10_generate_keys.conf
[Unit]
After=
After=network.target auditd.service sshdgenkeys.service
Wants=sshdgenkeys.service
